# deploy/opencodex.service
[Unit]
Description=openCodex Autonomous Service
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/openCodex
Environment=NODE_ENV=production
ExecStart=/usr/bin/env node %h/openCodex/src/index.js
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
